# supervisor/queue.ini
[program:queue]
command=php artisan queue:listen
process_name=%(program_name)s_%(process_num)02d
numprocs=3
directory=/var/www/backend
stdout_logfile=/var/log/%(program_name)s.log
autostart=true
autorestart=true
user=www-data
stopsignal=KILL
startsecs=3